# packer/cpu/scripts/buildx-builder.service
[Unit]
Description=Configure buildx to use BuildKit daemon
After=buildkitd.service
Requires=buildkitd.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=buildkite-agent
Group=buildkite-agent

# Wait for BuildKit socket to be available
ExecStartPre=/bin/bash -c 'for i in {1..30}; do [ -S /run/buildkit/buildkitd.sock ] && exit 0; sleep 1; done; exit 1'

# Create the buildx builder config pointing to buildkitd
ExecStart=/bin/bash -c '\
  docker buildx rm baked-vllm-builder 2>/dev/null || true; \
  docker buildx create \
    --name baked-vllm-builder \
    --driver remote \
    --use \
    unix:///run/buildkit/buildkitd.sock && \
  echo "buildx configured to use buildkitd"'

[Install]
WantedBy=multi-user.target
